# roles/kong/templates/kong_supervisord.conf.j2
[program:kong]
command= /opt/bin/start_kong.sh
autostart=true
autostart=true
autorestart=true
startretries=999
redirect_stderr=True
nodaemon = False
debug = False
loglevel = "info"
stopasgroup=true